# supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[unix_http_server]
file=/run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:manga-notify-bot]
command=/usr/bin/python3 -m manga_notify.main
autostart=true
autorestart=true
stderr_logfile=/var/log/manga_notify_bot.err.log
stdout_logfile=/var/log/manga_notify_bot.out.log

[program:manga-notify-jobs]
command=/usr/bin/python3 -m manga_notify.jobs
autostart=true
autorestart=true
stderr_logfile=/var/log/manga_notify_jobs.err.log
stdout_logfile=/var/log/manga_notify_jobs.out.log
